# blokas-jack/blokas-jack-1.0.2/debian/jack.service
[Unit]
Description=JACK Server
After=sound.target

[Service]
LimitRTPRIO=95
LimitMEMLOCK=infinity
Environment=JACK_NO_AUDIO_RESERVATION=1
Environment=JACK_PROMISCUOUS_SERVER=jack
ExecStart=/etc/jackdrc
User=jack
Group=jack
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
